# config/x2go.crontab
*/2 * * * *  /usr/bin/x2go-start-local.sh &> /dev/null

